# redis.conf
[program:redis]
command=/usr/bin/redis-server
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis-error.log
